# bot/scripts/vgen-farmer-okx.service
[Unit]
Description=vGen OKX volume farmer (v3 5M campaign — DEMO; switch --demo to --live + LIVE_OKX_ACK for real money)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/vGen-OKX/bot
EnvironmentFile=-/root/vGen-OKX/bot/.env
Environment="PYTHONUNBUFFERED=1"
# NOTE: paper mode accrues ZERO real volume — the 5M bounty needs --demo (rehearsal)
# or --live (real). For --live also set LIVE_OKX_ACK=I_UNDERSTAND in .env.
ExecStart=/root/vGen-OKX/bot/.venv/bin/python3 /root/vGen-OKX/bot/scripts/run_volume_farmer_okx.py \
    --config /root/vGen-OKX/bot/config/config_volume_farmer_okx_v3.yaml \
    --label okx-v3 \
    --demo \
    --duration-days 365
# on-failure, NOT always: a clean exit (volume target reached, operator stop)
# must STAY down — Restart=always looped startup/stop cards forever after the
# session halted on target (halted=true persists in the state file).
Restart=on-failure
RestartSec=15
StandardOutput=append:/root/vGen-OKX/bot/logs/farmer_okx.log
StandardError=append:/root/vGen-OKX/bot/logs/farmer_okx.log
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
