# bot/scripts/vgen-farmer-paper-1m.service
[Unit]
Description=vGen MEXC volume farmer (paper 1m)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/vGen/mexc-bot
EnvironmentFile=-/root/vGen/mexc-bot/.env
Environment="PYTHONUNBUFFERED=1"
ExecStart=/root/vGen/.venv/bin/python /root/vGen/mexc-bot/scripts/run_volume_farmer.py \
    --config /root/vGen/mexc-bot/config/config_volume_farmer_1m.yaml \
    --state-file volume_farmer_paper_1m_state.json \
    --label PAPER-1M \
    --poll-seconds 20 \
    --duration-days 365
Restart=always
RestartSec=10
StandardOutput=append:/root/vGen/mexc-bot/logs/farmer_paper_1m.log
StandardError=append:/root/vGen/mexc-bot/logs/farmer_paper_1m.log
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
